# smon.service
[Unit]
Description=server monitoring daemon

[Service]
Type=simple
ExecStart=/usr/bin/smon.py

[Install]
WantedBy=network.target
